# tools/ai.duckbot.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.duckbot.telegram</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/duckets/.openclaw/workspace/duck-cli-src/tools/run-duck-telegram.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/duckets/.openclaw/workspace/duck-cli-src</string>
    <key>StandardOutPath</key>
    <string>/tmp/duck-telegram.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/duck-telegram.err</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DUCK_CLI_PATH</key>
        <string>/Users/duckets/.openclaw/workspace/duck-cli-src/duck</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
